# conf/etc/services.d/01-exporter/run
#!/bin/sh
NAME="OpenDataSoft Exporter"
COMMAND="python3 /usr/bin/opendatasoft-exporter.py"

exec service-handler "$NAME" "$COMMAND" 
